# packaging/ulanzi-niri.service
[Unit]
Description=Ulanzi Stream Controller D200X driver (niri integration)
Documentation=https://github.com/doitian/ulanzi-studio-niri
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ulanzi-niri run
Restart=on-failure
RestartSec=3
TimeoutStopSec=5
KillMode=mixed
# This daemon's job is to spawn arbitrary user applications (kitty,
# firefox, niri msg, ...). Sandboxing flags such as ProtectHome= and
# ProtectSystem= are inherited by child processes via the mount
# namespace, which breaks any spawned app that needs to write to
# $HOME or /tmp (e.g. kitty's effective-config cache). Keep only
# NoNewPrivileges, which does not affect filesystem visibility.
NoNewPrivileges=true
Environment=ULANZI_NIRI_LOG=INFO

[Install]
WantedBy=graphical-session.target
